# publisher/generation/crontab
0 9,12,3 * * * python3 /generate_article.py